# infrastructure/systemd/github-runner.service
[Unit]
Description=GitHub Actions Runner for QuantumRishi
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=github-runner
Group=github-runner
WorkingDirectory=/opt/actions-runner
ExecStart=/opt/actions-runner/runsvc.sh
Restart=always
RestartSec=5
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=5min

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/actions-runner/_work
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
